# servers/mcow/homelab-admin.service
# mcow LXC privilege: N/A (KVM virtual machine, not LXC)   (per servers/mcow/lxc-probe.md, Plan 12-03)
# Hardening block selected: full strict (ProtectSystem=strict + PrivateTmp=yes + ProtectKernelTunables=yes)
# P-15 degraded path does NOT apply — mcow is KVM with full Ubuntu 24.04 kernel.

[Unit]
Description=Homelab Admin Dashboard (Next.js on Bun)
Documentation=https://github.com/makscee/homelab/tree/main/apps/admin
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=homelab-admin
Group=homelab-admin
WorkingDirectory=/opt/homelab-admin/current/homelab/apps/admin
EnvironmentFile=/etc/homelab-admin/env
Environment=NODE_ENV=production
Environment=PORT=3847
# bind on docker bridge gateway — compose Caddy reaches host services via
# host.docker.internal (172.17.0.1), not loopback (HMB-9). Host-only IP keeps
# :3847 off the public NIC since UFW is inactive on mcow.
Environment=HOSTNAME=172.17.0.1
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# --bun forces Bun runtime for child processes (next CLI spawns node by default;
# we need Bun globals so bun:sqlite native module resolves via the shim proxy).
ExecStart=/usr/local/bin/bun --bun run /opt/homelab-admin/current/homelab/apps/admin/server.js
Restart=on-failure
RestartSec=5s
# Resource limits — tune after 24h soak (Plan 10 + Phase 17 metrics)
LimitNOFILE=65536

# ── Hardening ────────────────────────────────────────────────────────────
# mcow is a KVM VM (systemd-detect-virt=kvm, kernel 6.8.0-107-generic).
# Full strict hardening is safe — no LXC user-namespace constraints apply.
# Block transcribed from lxc-probe.md "Active recommendation" (Plan 12-03).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ReadWritePaths=/opt/homelab-admin /var/lib/homelab-admin /tmp
Environment=AUDIT_DB_PATH=/var/lib/homelab-admin/audit.db

[Install]
WantedBy=multi-user.target
